# etc/init.d/powertop
#!/sbin/openrc-run

description="Extend the battery life of laptop by optimizing power settings with PowerTOP"

command="/usr/sbin/powertop"
command_args="--auto-tune"
command_background="no"

depend() {
    after bootmisc
}

start() {
    ebegin "Optimizing power settings with PowerTOP"
    start-stop-daemon --start --exec ${command} -- ${command_args}
    eend $?
}

# No need for a stop function as it's a oneshot operation
