# etc/init.d/powertop
#!/sbin/openrc-run

description="Optimizes power settings with PowerTOP's --auto-tune"

depend() {
	after local
}

start() {
	ebegin "Applying PowerTOP auto-tune settings"
	/usr/sbin/powertop --auto-tune
	eend $?
}
